# package/feeds/ltq_common/ifx-base-files-common/files/etc/init.d/ifx_show_version.sh
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org
# Copyright (C) 2007 infineon.com

START=99

bindir=/opt/ifx/bin

start() {
   # start the dsl daemon
   [ -e ${bindir}/show_version.sh ] && ${bindir}/show_version.sh
	[ -e /opt/ifx/image_version ] && cat /opt/ifx/image_version
}
