# package/FLTE_TestAllPeripherals/files/TestAllPeripherals.sh
#!/bin/sh /etc/rc.common
#
# Install led driver and make led device node

START=73


start() {
		TestAllPeripherals	
	}

