# package/feeds/ltq_feeds_netcomp_cpe/ifx_config_common_features/files-opt/nat/etc/init.d/port_trigger.sh
#!/bin/sh /etc/rc.common
#START=66
start() {
	        /etc/rc.d/init.d/port_trigger_init 2>/dev/null
}
