# package/FLTE_I2CApp/files/resetcodecconfig.sh
#!/bin/sh /etc/rc.common
#
# Install led driver and make led device node

START=74


start() {
              /usr/sbin/configurecodec
        }
