# package/feeds/ltq_feeds_netcomp_cpe/ifx_config_common_features/files-opt/wan/services/etc/init.d/dnrd.sh
#!/bin/sh /etc/rc.common
#START=36
start() {
	if [ -f /usr/sbin/dnrd ]; then
		rm -rf /ramdisk/etc/dnrd
		mkdir /ramdisk/etc/dnrd
		/usr/sbin/dnrd
	fi
}
