# package/feeds/ltq_feeds_netcomp_cpe/ifx_config_common_features/files/etc/init.d/node_gpio.sh
#!/bin/sh /etc/rc.common
#START=14
start() {
	/usr/sbin/mknod_util ifx_gpio /dev/ifx_gpio
}
