# package/feeds/ltq_feeds_netcomp_cpe/ifx_config_common_features/files/etc/init.d/resolv.sh
#!/bin/sh /etc/rc.common
#START=61
start() {
	if [ ! -s /etc/resolv.conf ]; then
		echo "nameserver 168.95.1.1" > /etc/resolv.conf
	fi
}
