# package/feeds/ltq_feeds_netcomp_cpe/ifx_config_common_features/files/etc/init.d/switch_api.sh
#!/bin/sh /etc/rc.common
#START=35
start() {
	/sbin/modprobe drv_switch_api
}
